# build/hack_scripts/hostfiles-hack.service
[Unit]
After=sshd.service

[Service]
ExecStart=/usr/local/bin/hostfiles-hack.sh

[Install]
WantedBy=default.target
